# cometai.service
[Unit]
Description=CometAI Production Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=cometai
Group=cometai
WorkingDirectory=/opt/cometai
ExecStart=/opt/cometai/venv/bin/python server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cometai

# Environment variables
Environment=ENVIRONMENT=production
Environment=HOST=0.0.0.0
Environment=PORT=8080
Environment=LOG_LEVEL=INFO
Environment=WORKERS=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cometai/logs
ReadWritePaths=/opt/cometai/.cache

# Resource limits
LimitNOFILE=65536
MemoryMax=20G

[Install]
WantedBy=multi-user.target
